# tron-lifecycle-manager/crontab
# START CRON JOB
6 * * * * /usr/bin/env /usr/local/bin/python /app/cleanup.py > /proc/1/fd/1 2>/proc/1/fd/2
# END CRON JOB
